# com.cliproxyapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cliproxyapi</string>
    <key>ProgramArguments</key>
    <array>
        <string>${HOME}/CLIProxyAPI/cli-proxy-api</string>
        <string>--config</string>
        <string>${HOME}/CLIProxyAPI/config.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/CLIProxyAPI/proxy.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/CLIProxyAPI/proxy-error.log</string>
    <key>WorkingDirectory</key>
    <string>${HOME}/CLIProxyAPI</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>${HOME}</string>
    </dict>
</dict>
</plist>
